# scripts/myapp-desktop.service
[Unit]
Description=MyApp Desktop
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/myapp-desktop
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
